# crontab
0 * * * * cd /app && ./app >> /dev/stdout 2>&1
